# templates/llm-webui-qwen3-8b/s6-overlay/s6-rc.d/openwebui/run
#!/command/with-contenv bash
set -e

echo "=== Starting Open WebUI service ==="

wait_for_vllm() {
    local max_attempts=30
    local attempt=1
    
    echo "Waiting for vLLM at $OPENAI_API_BASE_URL..."
    
    while [ $attempt -le $max_attempts ]; do
        if curl -s --connect-timeout 5 "$OPENAI_API_BASE_URL/models" > /dev/null 2>&1; then
            echo "✅ vLLM is ready!"
            return 0
        fi
        
        echo "Attempt $attempt/$max_attempts: waiting for vLLM..."
        sleep 10
        attempt=$((attempt + 1))
    done
    
    echo "❌ vLLM did not respond after $max_attempts attempts"
    echo "Continuing anyway..."
    return 0
}

wait_for_vllm

echo "Open WebUI configuration:"
echo "  Port: $PORT"
echo "  Data Directory: $DATA_DIR"
echo "  API Base URL: $OPENAI_API_BASE_URL"

mkdir -p "$DATA_DIR"

echo "Starting Open WebUI..."

exec open-webui serve --port "$PORT" --host "0.0.0.0"
